# openrc/fakessh
#!/sbin/openrc-run

supervisor=supervise-daemon

description="Fake SSH Server"
command="/usr/local/bin/fakessh"
command_args="-c /etc/fakessh/fakessh.toml"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

depend() {
	need net
	after firewall
}